# src/etc/init.d/teal
#!/bin/bash
# begin_generated_IBM_copyright_prolog
#
# This is an automatically generated copyright prolog.
# After initializing,  DO NOT MODIFY OR MOVE
# ================================================================
#
# (C) Copyright IBM Corp.  2010,2011
# Eclipse Public License (EPL)
#
# ================================================================
#
# end_generated_IBM_copyright_prolog
#
#	/etc/rc.d/init.d/teal
#
#       This service manages the IBM Toolkit for Event Analysis and Logging (TEAL) daemon
#
# chkconfig: 345 20 80
# pidfile: /var/lock/teal.pid
# 

# Source function library.
. /etc/init.d/functions

# Local functions and definitions
TEAL=/opt/teal/ibm/teal/teal.py
PIDFILE=/var/lock/teal.pid

start() {
	echo -n "Starting teal: "
	daemon --pidfile $PIDFILE $TEAL --daemon
	RC=$?
	echo
	[ $RC -eq 0 ] && touch /var/lock/subsys/teal
	return $RC
}	

stop() {
	echo -n "Shutting down teal: "
	killproc -p $PIDFILE $TEAL
	RC=$?
	echo
	[ $RC -eq 0 ] && rm -f /var/lock/subsys/teal
	return $RC
}

case "$1" in
    start)
	start
	;;
    stop)
	stop
	;;
    status)
	status -p $PIDFILE $TEAL
	;;
    restart|reload)
    stop
	start
	;;
    condrestart)
    	#<Restarts the servce if it is already running. For example:>
	[ -f /var/lock/subsys/<service> ] && restart || :
	;;
    *)
	echo "Usage: teal {start|stop|status|restart|reload}"
	exit 1
	;;
esac
exit $?
